# supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/sbin/init
priority=1

[program:sshd]
command=/usr/sbin/sshd -D -o UseDNS=no -o UsePAM=no
